# systemd/mediapi-app.service.template
[Unit]
Description=MediaPi Flask app
After=network-online.target mediapi-mpv.service
Wants=network-online.target mediapi-mpv.service

[Service]
Type=simple
User=${MEDIAPI_USER}
WorkingDirectory=${PROJECT_DIR}
ExecStart=${UV} run run.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
